# docker/s6-rc.d/migrations/up
#!/command/with-contenv sh
cd /app
echo "[migrations] Running Prisma migrations..."
exec s6-setuidgid node node node_modules/.bin/prisma migrate deploy
